# projects/toronto-alerts/toronto-alerts-flask/etc/crontab
SHELL=/bin/sh
0 * * * * echo "Hello from cron" >> /app/var/log/cron.log 2>&1
