# infrastructure/terraform/packer/scripts/cloudflared.service
[Unit]
Description=Cloudflare Tunnel (Vault ingress)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/cloudflared/env
ExecStart=/usr/bin/cloudflared --no-autoupdate tunnel run --token ${TUNNEL_TOKEN}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
